# crontab
# m h  dom mon dow   command
* * * * * cd ~/ddns && bash ~/ddns/cloudflare-wlan.sh
* * * * * cd ~/ddns && bash ~/ddns/cloudflare-eth.sh
* * * * * bash ~/direwolf-sdr.sh